# docker/s6/services.d/php-fpm/run
#!/usr/bin/execlineb -P
trap
{
    term
    {
        import -u ! kill -QUIT ${!}
    }
}
with-contenv
php-fpm
