# docker/s6/services.d/nginx/run
#!/usr/bin/execlineb -P
trap
{
    term
    {
        import -u ! kill -QUIT ${!}
    }
}
nginx -g "daemon off;"
